# uci-defaults/10_fix_passwall2.sh
sed -i '1c#!/bin/bash' /usr/share/passwall2/*.sh