# uci-defaults/10_socat_overlay_zh.sh
cat <<EOF | sh
mkdir -p /usr/lib/lua/luci/i18n
echo U29jYXQg5pivIExpbnV4IOS4i+eahOS4gOS4quWkmuWKn+iDveeahOe9kee7nOW3peWFt++8jOWQjeWtl+adpeeUseaYr+OAjFNvY2tldCBDQVTjgI3jgILlhbblip/og73kuI7mnInnkZ7lo6vlhpvliIDkuYvnp7DnmoQgTmV0Y2F0IOexu+S8vO+8jOWPr+S7peeci+WBmuaYryBOZXRjYXQg55qE5Yqg5by654mI44CCU29jYXQg6YWN572u54q25oCBAADlkK/nlKgAAOWkh+azqAAA5Y2P6K6uAADku4VJUHY2AOW9k+WLvumAieaXtu+8jOS7heebkeWQrElQdjbvvIzlkKbliJnlsIbkvJrlkIzml7bnm5HlkKxJUHY044CCAADnq6/lj6Povazlj5Hnm5HlkKzljY/orq7nm5HlkKznq6/lj6Pnu5HlrprliLDmnKzlnLDnq6/lj6MAAADnm67moIfljY/orq7nm67moIflnLDlnYDnm67moIfnq6/lj6PmiZPlvIDpmLLngavlopnnq6/lj6MAAAAP7gg5AAAAAQAAAUAAAAAMJ/KgewAAAAEAAADIAAAABjFiFNYAAAABAAABfAAAAAw1VZpvAAAAAQAAANgAAAAGQXXLuAAAAAEAAAFMAAAAFUlIpxgAAAABAAABiAAAABVc7enpAAAAAQAAAXAAAAAMXZrrHAAAAAEAAADgAAAAB2YzWckAAAABAAAAwAAAAAZ9G0b7AAAAAQAAASgAAAAMkVGKHgAAAAEAAAFkAAAADJXSJEgAAAABAAABNAAAAAy+8oqJAAAAAQAAALQAAAAMwnXsqgAAAAEAAADoAAAAPusdlt4AAAABAAAA0AAAAAb2VFCdAAAAAQAAAAAAAAC0AAABoA==|base64 -d > /usr/lib/lua/luci/i18n/socat.zh-cn.lmo
EOF
