# uci-defaults/91_info.sh
uci del dropbear.main.enable
uci del dropbear.main.Interface
uci del dropbear.main.RootPasswordAuth
uci commit dropbear
